# reports/lagerpro-cron.service
[Unit]
Description=LagerPro status runner

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/.openclaw/workspace
ExecStart=/home/ubuntu/.openclaw/workspace/reports/lagerpro-cron.sh
